# home/system/systemd/ollama.service.d/override.conf
[Service]
# Optimized for RTX 3070 8GB
Environment="OLLAMA_FLASH_ATTENTION=1"
Environment="OLLAMA_KV_CACHE_TYPE=q8_0"

# Upping context to 16k (Great for long code files)
Environment="OLLAMA_CONTEXT_LENGTH=16384"

# Pre-load the model so it stays in VRAM permanently
Environment="OLLAMA_KEEP_ALIVE=-1"

# Allow two parallel requests (e.g., chat + IDE autocomplete at once)
Environment="OLLAMA_NUM_PARALLEL=2"
