# home/system/systemd/displays-resume.service.tmpl
[Unit]
Description=Re-apply sway display layout after resume
DefaultDependencies=no
# ExecStop on this unit runs when leaving sleep (i.e. on wake).
Before=sleep.target
After=nvidia-resume.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/true
ExecStop={{ .chezmoi.homeDir }}/.local/bin/displays-resume
User={{ .chezmoi.username }}

[Install]
WantedBy=sleep.target
