# home/system/systemd/displays-resume.service.tmpl
[Unit]
Description=Re-apply sway display layout after resume
# Order after the sleep targets so this runs on resume, not at suspend entry.
# (systemd-suspend.service performs the actual sleep before these targets are
# reached, so units ordered After= them start once the machine wakes.)
After=suspend.target hibernate.target hybrid-sleep.target suspend-then-hibernate.target

[Service]
Type=oneshot
# Run as the desktop user so swaymsg/kanshi can reach the live sway session.
User={{ .chezmoi.username }}
ExecStart={{ .chezmoi.homeDir }}/.local/bin/displays-resume

[Install]
WantedBy=suspend.target hibernate.target hybrid-sleep.target suspend-then-hibernate.target
