# home/system/ssh/ssh.socket.d/port.conf
[Socket]
# Replace default :22 with custom port (must match sshd_config Port when using socket activation)
ListenStream=
ListenStream=0.0.0.0:9701
ListenStream=[::]:9701
