# home/dot_config/systemd/user/dropbox.service.tmpl
[Unit]
Description=Dropbox sync daemon
After=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
# dropbox-start derives DISPLAY from the XWayland socket so the tray icon
# appears; the systemd user manager has no DISPLAY at login (lazy XWayland).
ExecStart=%h/.local/bin/dropbox-start start
ExecStop=/usr/bin/python3 /usr/bin/dropbox stop

[Install]
WantedBy=graphical-session.target
